# infra/systemd/rostferrum.service
[Unit]
Description=RostFerrum FastAPI (Gunicorn/UvicornWorker)
After=network.target

[Service]
Type=notify
WorkingDirectory=/var/www/rostferrum
EnvironmentFile=/var/www/rostferrum/.env
ExecStart=/var/www/rostferrum/.venv/bin/gunicorn scripts.api:app \
  --workers 2 \
  --worker-class uvicorn.workers.UvicornWorker \
  --bind 127.0.0.1:8000 \
  --graceful-timeout 30 \
  --timeout 60
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target


